# dot_config/systemd/user/restic-backup-tool.timer
[Unit]
Description=Automated restic backups

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
