# dot_config/systemd/user/restic-backup-check.service.tmpl
[Unit]
Description=Automatic check of restic backup of $HOME
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStart={{ .chezmoi.homeDir }}/.local/bin/run-backup check
