# root/usr/lib/systemd/system/user@.service.d/20-fast-shutdown.conf
# Band-aid solution for apps like Bitwig that seem to ignore SIGTERM

[Service]
TimeoutStopSec=10s
