# root/etc/systemd/system/mnt-win.automount
[Unit]
Description=Automount /mnt/win

[Automount]
Where=/mnt/win

[Install]
WantedBy=multi-user.target
