# root/etc/systemd/system/srv-oldbackup.automount
[Unit]
Description=Automount /srv/oldbackup

[Automount]
Where=/srv/oldbackup
TimeoutIdleSec=300

[Install]
WantedBy=multi-user.target
