# root/etc/systemd/system/srv-public.automount
[Unit]
Description=Automount /srv/public

[Automount]
Where=/srv/public
TimeoutIdleSec=300

[Install]
WantedBy=multi-user.target
