# root/etc/systemd/system/srv-backup.automount
[Unit]
Description=Automount /srv/backup
After=network.target network-online.target remote-fs-pre.target nss-lookup.target

[Automount]
Where=/srv/backup
TimeoutIdleSec=300

[Install]
WantedBy=multi-user.target
