# root/etc/systemd/system/srv-public.automount
[Unit]
Description=Automount /srv/public
After=network.target network-online.target remote-fs-pre.target nss-lookup.target

[Automount]
Where=/srv/public
TimeoutIdleSec=300

[Install]
WantedBy=multi-user.target
